# quadlets/query-log-maintenance.service
# ~/.config/systemd/user/query-log-maintenance.service
# Runs partition maintenance for query_log: creates upcoming partitions,
# drops expired ones. Triggered daily by query-log-maintenance.timer.

[Unit]
Description=Query log partition maintenance
After=postgres.service

[Service]
Type=oneshot
EnvironmentFile=%h/.config/llm-stack/env
# Credentials sourced from ragstack.env — provides DATABASE_URL
EnvironmentFile=%h/.config/llm-stack/ragstack.env
ExecStart=/bin/sh %h/git/rag-suite/migrations/maintain-partitions.sh
